# controller/DEBIAN/etc/init/controller.conf
#
# Palette controller
#
description "Palette controller daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/controller /etc/controller.ini
