# controller/DEBIAN/debian/controller.postrm
#!/bin/sh

update-rc.d controller remove

exit 0
