# controller/DEBIAN/debian/controller.postrm
#!/bin/sh

update-rc.d -f controller remove

exit 0
